# cron/cron.daily/asolearning
#!/bin/sh

#CSVふぁいる削除
sh /data/batch/csvFileDelete.sh

#DBバックアップ
sh /data/batch/dbbackup.sh

exit 0
